# scripts/cleanup-dblck.service
[Unit]
Description=cleanup stale devtools db.lck files

[Service]
Type=exec
ExecStart=cleanup-dblck
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
